# super_turtle/claude-telegram-bot/launchagent/com.claude-telegram-ts.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude-telegram-ts</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>cd /Users/USERNAME/Dev/agentic &amp;&amp; exec "${SUPERTURTLE_BIN:-superturtle}" service run</string>
    </array>

    <key>WorkingDirectory</key>
    <!-- Update this path to your repo root -->
    <string>/Users/USERNAME/Dev/agentic</string>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- Required: Your Telegram bot token -->
        <key>TELEGRAM_BOT_TOKEN</key>
        <string>your-bot-token-here</string>

        <!-- Required: Comma-separated Telegram user IDs -->
        <key>TELEGRAM_ALLOWED_USERS</key>
        <string>123456789</string>

        <!-- Recommended: Working directory for Claude -->
        <key>CLAUDE_WORKING_DIR</key>
        <string>/Users/USERNAME/Dev/agentic</string>

        <!-- Optional: Global npm binary or source checkout launcher -->
        <key>SUPERTURTLE_BIN</key>
        <string>superturtle</string>

        <!-- Optional: For voice transcription -->
        <key>OPENAI_API_KEY</key>
        <string>sk-...</string>

        <!-- Optional: Comma-separated paths Claude can access -->
        <key>ALLOWED_PATHS</key>
        <string>/Users/USERNAME/Dev,/Users/USERNAME/Documents</string>

        <!-- Optional: Rate limiting -->
        <key>RATE_LIMIT_ENABLED</key>
        <string>true</string>
        <key>RATE_LIMIT_REQUESTS</key>
        <string>40</string>
        <key>RATE_LIMIT_WINDOW</key>
        <string>60</string>

        <!-- Ensure PATH includes necessary binaries -->
        <key>PATH</key>
        <string>/Users/USERNAME/.bun/bin:/Users/USERNAME/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <!-- For source checkouts without the npm package, replace SUPERTURTLE_BIN with:
         /Users/USERNAME/Dev/agentic/super_turtle/bin/superturtle.js -->

    <!-- Restart automatically if it crashes -->
    <key>KeepAlive</key>
    <true/>

    <!-- Start on login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Log output (replace TOKEN_PREFIX with your bot token's numeric prefix, e.g. 123456) -->
    <key>StandardOutPath</key>
    <string>/tmp/claude-telegram-bot.TOKEN_PREFIX.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/claude-telegram-bot.TOKEN_PREFIX.err</string>
</dict>
</plist>
